# .docker/crontab
*/5 * * * * /usr/unrealircd/unrealircd croncheck
